# supervisor/supervisor-npm.conf
[program:npm]
command=npm run prod
directory=/var/www/html
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr