# supervisor/supervisor-npm.conf
[program:npm]
command=npm run dev
directory=/var/www/html