# supervisor/supervisor-npm.conf
[program:npm]
command=npm run dev
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/npm_stdout.log
stderr_logfile=/var/log/supervisor/npm_stderr.log
