# supervisor/supervisor-npm.conf
[program:npm]
command=npm run prod
directory=/var/www/html